# deploy/systemd/support-agent-frontend-experimental.service
[Unit]
Description=Support Agent Frontend experimental bridge
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/support-agent-frontend/bridge
EnvironmentFile=-/etc/default/support-agent-frontend-experimental
ExecStart=/usr/bin/env node /opt/support-agent-frontend/bridge/src/server.js
Restart=always
RestartSec=3
User=ubuntu
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
