# deploy/systemd/support-chat-runtime-experimental.service
[Unit]
Description=Support Chat Runtime experimental bridge
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/support-chat-runtime/bridge
EnvironmentFile=-/etc/default/support-chat-runtime-experimental
ExecStart=/usr/bin/env node /opt/support-chat-runtime/bridge/src/server.js
Restart=always
RestartSec=3
User=ubuntu
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
